# conf/sys/supervisord/supervisord.conf
[program:testapp]
user={{user}}
command={{virtualenv_dir}}/bin/uwsgi -i {{app_dir}}/conf/uwsgi/uwsgi.ini
directory={{virtualenv_dir}}/src
autostart=true
autorestart=true
stopsignal=INT
